# watcher/systemd/model-connector-watcher.service
[Unit]
Description=Model Connector watcher
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
User=model-connector-watcher
Group=model-connector-watcher
EnvironmentFile=/etc/model-connector-watcher/watcher.env
ExecStart=/usr/bin/node /opt/model-connector-watcher/current/dist/src/main.js --config /etc/model-connector-watcher/config.yaml
WatchdogSec=90s
Restart=on-failure
RestartSec=10s
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/model-connector-watcher
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
